# _info/service/airodump.service
[Unit]
Description=airodump

[Service]
Type=idle
ExecStart=/bin/bash /home/pi/airodump.sh
